# src/sentry/postman/drkonqi-sentry-postman.path
# SPDX-License-Identifier: GPL-2.0-only OR GPL-3.0-only OR LicenseRef-KDE-Accepted-GPL
# SPDX-FileCopyrightText: 2023 Harald Sitter <sitter@kde.org>

[Unit]
Description=Submitting pending crash events

[Path]
PathExistsGlob=%C/drkonqi/sentry-envelopes/*
MakeDirectory=true
TriggerLimitIntervalSec=30 minutes
TriggerLimitBurst=4

[Install]
WantedBy=default.target plasma-core.target
